# S2/etc/systemd/system/xenos-control-dns-2.service
[Unit]
Description=xenos-control-dns-2.sh

[Service]
Type=oneshot
ExecStart=/usr/bin/xenos-control-dns-2.sh

[Install]
WantedBy=multi-user.target